# deploy/systemd/agentdeck-app-server@.service
[Unit]
Description=Agent Deck Codex app-server instance %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agentdeck
Group=agentdeck
WorkingDirectory=/opt/agentdeck
Environment=HOME=/var/lib/agentdeck/home
EnvironmentFile=/etc/agentdeck/agentdeck-app-server-%i.env
ExecStart=/usr/local/bin/codex app-server --listen ${CODEX_APP_SERVER_LISTEN} -c approval_policy=\"never\" -c sandbox_mode=\"danger-full-access\"
Restart=always
RestartSec=2
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/agentdeck /opt/agentdeck /tmp /run

[Install]
WantedBy=multi-user.target
