# deploy/systemd/agentdeck-app-server@.service
[Unit]
Description=Agent Deck Codex app-server instance %i
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=@AGENTDECK_RUN_USER@
Group=@AGENTDECK_RUN_GROUP@
WorkingDirectory=@AGENTDECK_CURRENT_DIR@
Environment=HOME=@AGENTDECK_HOME@
EnvironmentFile=@AGENTDECK_ENV_DIR@/agentdeck-app-server-%i.env
ExecStart=@CODEX_BIN@ app-server --listen ${CODEX_APP_SERVER_LISTEN} -c approval_policy=\"@CODEX_APPROVAL_POLICY@\" -c sandbox_mode=\"@CODEX_SANDBOX_MODE@\"
Restart=on-failure
RestartSec=5
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=@AGENTDECK_DATA_DIR@ @AGENTDECK_CURRENT_DIR@ /tmp /run

[Install]
WantedBy=multi-user.target
