# deploy/systemd/agentdeck-app-server@.service
[Unit]
Description=Agent Deck Codex app-server instance %i
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/stacks/agentdeck
Environment=HOME=/home/ubuntu
EnvironmentFile=@AGENTDECK_ENV_DIR@/agentdeck-app-server-%i.env
ExecStart=/usr/local/bin/codex app-server --listen ${CODEX_APP_SERVER_LISTEN} -c approval_policy=\"never\" -c sandbox_mode=\"danger-full-access\"
Restart=on-failure
RestartSec=5
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/data/agentdeck /opt/stacks/agentdeck /tmp /run

[Install]
WantedBy=multi-user.target
