# deploy/systemd/agentdeck-runtime.service
[Unit]
Description=Agent Deck persistent agent runtime
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=@AGENTDECK_RUN_USER@
Group=@AGENTDECK_RUN_GROUP@
WorkingDirectory=@AGENTDECK_CURRENT_DIR@-runtime
Environment=HOME=@AGENTDECK_HOME@
Environment=DATA_DIR=@AGENTDECK_DATA_DIR@
Environment=CODEX_BIN=@CODEX_BIN@
Environment=PATH=@AGENTDECK_DATA_DIR@/provider-tools/bin:@AGENTDECK_HOME@/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PLAYWRIGHT_BROWSERS_PATH=@AGENTDECK_DATA_DIR@/cache/ms-playwright
Environment=RUNTIME_HOST=127.0.0.1
Environment=RUNTIME_PORT=3852
Environment=RUNTIME_MODE=active
Environment=AGENTDECK_SYSTEMD_UNIT_VERSION=2
Environment=CODEX_APP_SERVER_PORT_BASE=4620
EnvironmentFile=@AGENTDECK_ENV_DIR@/runtime.env
ExecStart=/usr/bin/flock -n /run/agentdeck/runtime-production.lock /usr/bin/node @AGENTDECK_CURRENT_DIR@-runtime/server/dist/agentdeck-runtime.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=660
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=@AGENTDECK_DATA_DIR@ /tmp /run /run/agentdeck

[Install]
WantedBy=multi-user.target
