# deploy/systemd/agentdeck-web.service
[Unit]
Description=Agent Deck web gateway
After=network-online.target agentdeck-runtime.service
Wants=network-online.target

[Service]
Type=simple
User=agentdeck
Group=agentdeck
WorkingDirectory=/opt/agentdeck
Environment=HOME=/var/lib/agentdeck/home
Environment=USE_AGENT_RUNTIME=1
Environment=AGENT_RUNTIME_URL=http://127.0.0.1:3852
EnvironmentFile=/etc/agentdeck/web.env
ExecStart=/usr/bin/node /opt/agentdeck/server/dist/index.js
Restart=always
RestartSec=3
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/agentdeck /opt/agentdeck /tmp

[Install]
WantedBy=multi-user.target
