# deploy/systemd/agentdeck-runtime.service
[Unit]
Description=Agent Deck persistent agent runtime
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/stacks/agentdeck/current
Environment=HOME=/home/ubuntu
Environment=DATA_DIR=/opt/data/agentdeck
Environment=PATH=/opt/data/agentdeck/provider-tools/bin:/usr/local/bin:/usr/bin:/bin:/home/ubuntu/.local/bin
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/data/agentdeck/cache/ms-playwright
Environment=RUNTIME_HOST=127.0.0.1
Environment=RUNTIME_PORT=3852
Environment=RUNTIME_MODE=active
Environment=CODEX_APP_SERVER_PORT_BASE=4620
EnvironmentFile=@AGENTDECK_ENV_DIR@/runtime.env
ExecStart=/usr/bin/flock -n /run/agentdeck/runtime-production.lock /usr/bin/node /opt/stacks/agentdeck/current/server/dist/agentdeck-runtime.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/data/agentdeck /opt/stacks/agentdeck /tmp /run /run/agentdeck

[Install]
WantedBy=multi-user.target
