# deploy/systemd/agentdeck-web.service
[Unit]
Description=Agent Deck web gateway
After=network-online.target agentdeck-runtime.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=@AGENTDECK_RUN_USER@
Group=@AGENTDECK_RUN_GROUP@
WorkingDirectory=@AGENTDECK_CURRENT_DIR@-web
Environment=HOME=@AGENTDECK_HOME@
Environment=DATA_DIR=@AGENTDECK_DATA_DIR@
Environment=CODEX_BIN=@CODEX_BIN@
Environment=PATH=@AGENTDECK_DATA_DIR@/provider-tools/bin:@AGENTDECK_HOME@/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PLAYWRIGHT_BROWSERS_PATH=@AGENTDECK_DATA_DIR@/cache/ms-playwright
Environment=USE_AGENT_RUNTIME=1
Environment=AGENT_RUNTIME_URL=http://127.0.0.1:3852
EnvironmentFile=@AGENTDECK_ENV_DIR@/web.env
ExecStart=/usr/bin/node @AGENTDECK_CURRENT_DIR@-web/server/dist/index.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=@AGENTDECK_DATA_DIR@ /tmp /run

[Install]
WantedBy=multi-user.target
