# deploy/systemd/codex-app-server@.service
[Unit]
Description=Agent Deck Codex app-server instance %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/stacks/codex-mobile
Environment=HOME=/home/ubuntu
EnvironmentFile=/opt/data/codex-mobile/codex-app-server-%i.env
ExecStart=/usr/local/bin/codex app-server --listen ${CODEX_APP_SERVER_LISTEN} -c approval_policy=\"never\" -c sandbox_mode=\"danger-full-access\"
Restart=always
RestartSec=2
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/data/codex-mobile /opt/stacks/codex-mobile /home/ubuntu /tmp /run

[Install]
WantedBy=multi-user.target
